# _opt_gitlab/gitlab/sv/nginx/run
#!/bin/sh
exec 2>&1

cd /var/opt/gitlab/nginx
exec chpst -P /opt/gitlab/embedded/sbin/nginx -p /var/opt/gitlab/nginx
