# _opt_gitlab/gitlab/sv/gitlab-exporter/run
#!/bin/sh
exec 2>&1

umask 077
exec chpst -P \
  -U git:git \
  -u git:git \
  /opt/gitlab/embedded/bin/gitlab-exporter web \
    -c /var/opt/gitlab/gitlab-exporter/gitlab-exporter.yml
