# _opt_gitlab/gitlab/sv/redis/run
#!/bin/sh
exec 2>&1

umask 077
exec chpst -P -U gitlab-redis:gitlab-redis -u gitlab-redis:gitlab-redis /opt/gitlab/embedded/bin/redis-server /var/opt/gitlab/redis/redis.conf
